# docker/supervisor/supervisord.conf
[supervisord]
# Mantiene el proceso principal en primer plano para que el contenedor siga activo.
nodaemon=true
user=root

[program:php-fpm]
# Ejecuta PHP para procesar las peticiones de Laravel.
command=/usr/local/sbin/php-fpm -F
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
# Servidor web que recibe las peticiones HTTP en el puerto 3000.
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
